# supervisor.conf
[program:main]
command=/bin/bash /Fooocus-API/check.sh
directory=/Fooocus-API/
autostart=true
autorestart=true
startsecs = 3
stderr_logfile=/main.err.log
stdout_logfile=/main.out.log
